# config/supervisord.conf
# 🦅 NoC Raven - Supervisor Configuration
# Process management for containerized telemetry appliance

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=nocraven:nocraven

[inet_http_server]
port=127.0.0.1:9001
username=admin
password=noc-raven-admin

[supervisord]
logfile=/var/log/noc-raven/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
umask=022
user=root
identifier=supervisor
directory=/tmp
nocleanup=true
childlogdir=/var/log/noc-raven
strip_ansi=false
silent=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

# =============================================================================
# CORE TELEMETRY SERVICES
# =============================================================================

[program:goflow2]
command=/opt/noc-raven/bin/goflow2 -config /opt/noc-raven/config/goflow2.yml
directory=/opt/noc-raven
user=nocraven
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
redirect_stderr=false
stdout_logfile=/var/log/noc-raven/goflow2.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
stderr_logfile=/var/log/noc-raven/goflow2.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stderr_capture_maxbytes=1MB
stderr_events_enabled=false
environment=NOC_RAVEN_HOME="/opt/noc-raven",DATA_PATH="/data",USER="nocraven",HOME="/opt/noc-raven"
serverurl=AUTO
priority=100

[program:fluent-bit]
command=/usr/bin/fluent-bit --config /etc/fluent-bit/fluent-bit.conf
directory=/opt/noc-raven
user=nocraven
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
redirect_stderr=false
stdout_logfile=/var/log/noc-raven/fluent-bit.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/noc-raven/fluent-bit.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=NOC_RAVEN_HOME="/opt/noc-raven",DATA_PATH="/data",USER="nocraven"
priority=200

[program:vector]
command=/usr/bin/vector --config /etc/vector/vector.toml
directory=/opt/noc-raven
user=nocraven
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=15
stopasgroup=false
killasgroup=false
redirect_stderr=false
stdout_logfile=/var/log/noc-raven/vector.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/noc-raven/vector.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=NOC_RAVEN_HOME="/opt/noc-raven",DATA_PATH="/data",USER="nocraven"
priority=300

[program:telegraf]
command=/usr/bin/telegraf --config /etc/telegraf/telegraf.conf
directory=/opt/noc-raven
user=nocraven
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
redirect_stderr=false
stdout_logfile=/var/log/noc-raven/telegraf.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/noc-raven/telegraf.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=NOC_RAVEN_HOME="/opt/noc-raven",DATA_PATH="/data",USER="nocraven"
priority=400

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
directory=/opt/noc-raven
user=root
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
exitcodes=0
stopsignal=QUIT
stopwaitsecs=30
stopasgroup=false
killasgroup=false
redirect_stderr=false
stdout_logfile=/var/log/noc-raven/nginx.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/noc-raven/nginx.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
priority=500

# =============================================================================
# MONITORING AND MANAGEMENT
# =============================================================================

[program:health-monitor]
command=/opt/noc-raven/bin/health-check.sh --quiet
directory=/opt/noc-raven
user=nocraven
numprocs=1
autostart=false
autorestart=false
startsecs=0
startretries=0
exitcodes=0,1,2
stopsignal=TERM
stopwaitsecs=5
stdout_logfile=/var/log/noc-raven/health-monitor.log
stderr_logfile=/var/log/noc-raven/health-monitor.log

# Log rotation service
[program:logrotate]
command=/usr/sbin/logrotate -f /etc/logrotate.d/noc-raven
directory=/
user=root
numprocs=1
autostart=false
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/var/log/noc-raven/logrotate.log
stderr_logfile=/var/log/noc-raven/logrotate.log

# =============================================================================
# EVENT LISTENERS
# =============================================================================

[eventlistener:crashmail]
command=/opt/noc-raven/scripts/crash-notifier.sh
events=PROCESS_STATE_FATAL
buffer_size=100
stdout_logfile=/var/log/noc-raven/crashmail.log
stderr_logfile=/var/log/noc-raven/crashmail.log

[eventlistener:memmon]
command=memmon -a 400MB -m support@rectitude369.com
events=TICK_60
buffer_size=10
stdout_logfile=/var/log/noc-raven/memmon.log
stderr_logfile=/var/log/noc-raven/memmon.log

# =============================================================================
# GROUPS
# =============================================================================

[group:telemetry]
programs=goflow2,fluent-bit,vector,telegraf
priority=100

[group:web]
programs=nginx
priority=200

[group:monitoring]
programs=health-monitor,logrotate
priority=300

# =============================================================================
# FCGI PROGRAMS (if needed for web interface extensions)
# =============================================================================

# [fcgi-program:web-api]
# socket=tcp://localhost:9000
# command=/opt/noc-raven/bin/web-api-server
# directory=/opt/noc-raven
# user=nocraven
# numprocs=4
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/noc-raven/web-api.log
# stderr_logfile=/var/log/noc-raven/web-api.log
